# quandoo/roles/docker/templates/supervisord.conf.j2
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
files=/etc/ssh/ssh_config
